# system/systemd/dashboard-stats.timer
[Unit]
Description=Trigger dashboard metrics generation every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2m
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
